# images/test-moodle/assets/crontab
* * * * * root /usr/bin/php /var/www/html/moodle/admin/cli/cron.php > /dev/null 2> /opt/moodledata/cron_error.log
 